# grafana/rootfs/etc/services.d/grafana/finish
#!/usr/bin/with-contenv bash
# Give s6 a moment before exit; harmless placeholder
sleep 1
