# grafana/rootfs/etc/services.d/grafana/run
#!/usr/bin/with-contenv bash
set -euo pipefail

if getent passwd grafana >/dev/null 2>&1; then
  chown -R grafana:grafana /data || true
  exec s6-setuidgid grafana grafana-server \
    --homepath=/usr/share/grafana \
    --config=/etc/grafana/grafana.ini \
    --packaging=addon \
    --pidfile=/run/grafana.pid
else
  # Fallback: run as root (shouldn't happen after Dockerfile fix)
  exec grafana-server \
    --homepath=/usr/share/grafana \
    --config=/etc/grafana/grafana.ini \
    --packaging=addon \
    --pidfile=/run/grafana.pid
fi
