# grafana/rootfs/etc/services.d/grafana/run
#!/usr/bin/with-contenv bash
set -euo pipefail

OPTS=/data/options.json
ADMIN_USER=$(jq -r '.admin_user' "$OPTS")
ADMIN_PASSWORD=$(jq -r '.admin_password' "$OPTS")
ALLOW_ANON=$(jq -r '.allow_anonymous' "$OPTS")

# Query Supervisor (no bashio)
INFO="$(curl -fsS -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" http://supervisor/addons/self/info)"
SELF_SLUG="$(echo "$INFO" | jq -r '.data.slug')"                     # e.g. 5dcfe9b4_energymanager_grafana
DIRECT_PORT="$(echo "$INFO" | jq -r '.data.ports["3000/tcp"] // empty')"

# Decide root_url: direct port -> "/", otherwise HA slug path (never token)
unset GF_SERVER_ROOT_URL || true
if [ -n "${DIRECT_PORT}" ] && [ "${DIRECT_PORT}" != "null" ]; then
  export GF_SERVER_ROOT_URL="/"
else
  export GF_SERVER_ROOT_URL="/hassio/ingress/${SELF_SLUG}/"
fi

# Ingress/subpath
export GF_SERVER_SERVE_FROM_SUB_PATH=true
export GF_SECURITY_ALLOW_EMBEDDING=true
export GF_SERVER_ENFORCE_DOMAIN=false
export GF_AUTH_LOGIN_COOKIE_NAME=grafana_ingress

# Writable paths
export GF_PATHS_DATA=/data/grafana
export GF_PATHS_LOGS=/data/logs
export GF_PATHS_PLUGINS=/data/plugins
export GF_PATHS_PROVISIONING=/etc/grafana/provisioning

# Admin / anonymous
export GF_SECURITY_ADMIN_USER="$ADMIN_USER"
export GF_SECURITY_ADMIN_PASSWORD="$ADMIN_PASSWORD"
if [ "$ALLOW_ANON" = "true" ]; then
  export GF_AUTH_ANONYMOUS_ENABLED=true
  export GF_AUTH_ANONYMOUS_ORG_ROLE=Viewer
else
  export GF_AUTH_ANONYMOUS_ENABLED=false
fi

# Start
if getent passwd grafana >/dev/null 2>&1; then
  chown -R grafana:grafana /data || true
  exec s6-setuidgid grafana grafana-server \
    --homepath=/usr/share/grafana \
    --config=/etc/grafana/grafana.ini \
    --packaging=addon
else
  exec grafana-server \
    --homepath=/usr/share/grafana \
    --config=/etc/grafana/grafana.ini \
    --packaging=addon
fi

